# deploy/astroarch-bridge.service
[Unit]
Description=Astroarch Bridge - Backend daemon for Astroarch Interface mobile app
Documentation=https://github.com/zarletti/astroarch-interface
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python -m astroarch_bridge
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=%h/.config/astroarch-bridge %h/Pictures/Ekos

# IMPORTANT: NOTHING PrivateTmp=true. The bridge must be able to read
# the xauth file of the user logged in graphically (e.g., /tmp/xauth_XXXXXX
# created by sddm) to spawn KStars/PHD2 on its display. With
# PrivateTmp=true, the service sees an isolated /tmp and the xauth files
# become invisible → "could not connect to display :0".
LockPersonality=true
RestrictRealtime=true

# Environment overrides (optional - uncomment to customize)
# Environment=ASTROARCH_PORT=8765
# Environment=ASTROARCH_INDI_HOST=127.0.0.1
# Environment=ASTROARCH_INDI_PORT=7624
# Environment=ASTROARCH_PHD2_HOST=127.0.0.1
# Environment=ASTROARCH_PHD2_PORT=4400
# Environment=ASTROARCH_LOG_LEVEL=INFO

[Install]
WantedBy=default.target
